# app/systemd/logpipe.service
[Unit]
Description=logpipe: stream FIFO writes to OpenSearch
Documentation=https://github.com/nettrash/logpipe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Pre-create the FIFO as root so the daemon (running as `logpipe`) can
# open it even when the path lives under /dev. The leading `+` runs this
# step with full privileges, regardless of User=.
ExecStartPre=+/bin/sh -c '\
  if [ ! -p /dev/logpipe ]; then \
    mkfifo -m 0622 /dev/logpipe; \
  fi; \
  chown logpipe:logpipe /dev/logpipe; \
  chmod 0622 /dev/logpipe'

ExecStart=/usr/local/bin/logpipe --config /etc/logpipe/config.toml

ExecStopPost=+/bin/sh -c 'rm -f /dev/logpipe'

Restart=on-failure
RestartSec=5s

# Bound shutdown. The daemon flushes within ~10s of SIGTERM; if it ever
# wedges, fall back to SIGKILL instead of waiting out the 90s default.
TimeoutStopSec=15s

# Run as a dedicated unprivileged user. Create with:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin logpipe
User=logpipe
Group=logpipe

# State / config directories (not strictly required for the FIFO model,
# but useful if you later add buffering on disk).
StateDirectory=logpipe
StateDirectoryMode=0750
ConfigurationDirectory=logpipe
ConfigurationDirectoryMode=0750

# --- Sandboxing -----------------------------------------------------------
# NOTE: PrivateDevices=true would hide /dev/logpipe from the rest of the
# system, defeating the whole point. If you relocate the FIFO under /run
# you can flip PrivateDevices back on.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
